# luci/luci-theme-prpl/root/etc/uci-defaults/30_luci-theme-prpl
#!/bin/sh

if [ "$PKG_UPGRADE" != 1 ]; then
	uci batch <<-EOF
		set luci.themes.PrplWrt=/luci-static/prpl
		set luci.main.mediaurlbase=/luci-static/prpl
		commit luci
	EOF
fi

exit 0
